# supervisorctl.conf
[program:vehicle_genius_portal]
command=bash start.sh
directory=/home/dimo/vehicle-genius-portal
